# upload-on-shutdown.service
[Unit]
Description=Upload Application Log to S3 on Shutdown/Reboot (${STAGE})
# Ensure the mount point for logs is available before running this service
RequiresMountsFor=/root/springlog
DefaultDependencies=no
After=network-online.target
Before=shutdown.target reboot.target umount.target final.target

[Service]
Type=oneshot
# Keep the service state as 'active' (though dead) after execution, for easier debugging with systemctl status
RemainAfterExit=yes
EnvironmentFile=/etc/default/upload_on_shutdown_env
# Environment variables for the script.
# Terraform will substitute ${S3_BUCKET_NAME} and ${STAGE} with the actual values.
Environment="S3_BUCKET_NAME=${S3_BUCKET_NAME}"
Environment="LOG_DIR_HOST=/root/springlog"
Environment="STAGE=${STAGE}"

ExecStart=/usr/local/bin/upload_on_shutdown.sh
# Give the script more time to complete, e.g., 90 seconds
TimeoutSec=90
User=root
Group=root
# Direct output to both journal and console for maximum visibility
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=shutdown.target reboot.target
